# autostart/open_html.desktop
[Desktop Entry]
Name=Attendence
Exec=#/home/home/github/Kids_Vids/open_html.py
Type=Application
X-GNOME-Autostart-enabled=true
