# autostart/run_django_server.desktop
[Desktop Entry]
Name=Attendence
Exec=/home/home/github/Kids_Vids/autostart/run_django_server.sh
Type=Application
X-GNOME-Autostart-enabled=true
